# docker/s6/services/imgproxy.svc/run
#!/bin/sh

echo "Launching imgdeflator service..."
/imgdeflator/imgdeflator
