# docker/s6/services/imgproxy.svc/run
#!/bin/sh

echo "Launching imgproxy service..."
/imgproxy/imgproxy
